# payload/Library/LaunchDaemons/com.erikng.installapplications.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.erikng.installapplications</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python</string>
		<string>/Library/Application Support/installapplications/installapplications.py</string>
		<string>--jsonurl</string>
		<string>https://tgroupmail.com/bootstrap.json</string>
		<string>--depnotify</string>
		<string>Status: Now setting up and configuring your machine...</string>
  		<string>Command: Image: /Library/Application Support/installapplications/cclogo.jpg</string>
		<string>Command: NotificationOn:</string>
		<string>Command: MainText: \nWe are setting up your Mac with a standard suite of software, configuration, and security policies. This process should only take a few minutes, so please don't restart your computer until we are finished.\n\nAt the end you will be asked to logout to complete process.\n</string>
		<string>Command: Logout: Device Enrollment complete. Please logout to begin FileVault encryption.</string>
		<string>DEPNotifyPath: /Applications/Utilities/DEPNotify.app</string>
		<string>DEPNotifyArguments: -munki -fullScreen</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/installapplications.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/installapplications.log</string>
</dict>
</plist>
